# debian/firesquid.service
[Unit]
AssertPathExists=/usr/bin/firesquid

[Service]
WorkingDirectory=~
ExecStart=/usr/bin/firesquid
Restart=always
PrivateTmp=true
NoNewPrivileges=true

[Install]
Alias=firesquid
WantedBy=default.target